# build/jenkins/slave-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--

http://en.wikipedia.org/wiki/Launchd
http://sourceforge.net/projects/lingon/files/
http://developer.apple.com/library/mac/#documentation/MacOSX/Conceptual/BPSystemStartup/Chapters/CreatingLaunchdJobs.html

copy this file into /Library/LaunchDaemons/

test:
/Library/LaunchDaemons/
sudo launchctl load -w slave-agent.plist
sudo launchctl unload -w slave-agent.plist

-->
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>slave-agent</string>
		<key>ProgramArguments</key>
		<array>
			<string>/var/jenkins/slave-agent.sh</string>
			<string>think</string>
		</array>
		<key>RunAtLoad</key>
		<true/>
	</dict>
</plist>
